# packaging/systemd/acva.target
[Unit]
Description=acva voice stack (orchestrator + backends)
# M4B post-consolidation: a single Speaches unit replaces the
# pre-M4B acva-whisper + acva-piper pair.
Wants=acva-llama.service acva-speaches.service acva.service
After=acva-llama.service acva-speaches.service

[Install]
WantedBy=default.target
