# packaging/systemd/acva.service
[Unit]
Description=acva — Autonomous Conversational Voice Agent
# M4B: backends collapsed to llama + speaches.
After=acva-llama.service acva-speaches.service
Wants=acva-llama.service acva-speaches.service

[Service]
Type=simple
ExecStart=%h/.local/bin/acva --config %h/.config/acva/config.yaml
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal
# Real-time scheduling for the audio thread (optional, but reduces underruns).
LimitRTPRIO=20
LimitMEMLOCK=infinity

[Install]
WantedBy=default.target
